# crontab
* * * * * . $HOME/.profile; $VARNISH_UPDATE
