# crontab
* * * * * /bin/sh /root/varnishconfig/update.sh
* * * * * /root/varnishconfig/doecho.sh
