# crontab
* * * * * root /root/varnishconfig/update.sh
