# crontab
0 0 * * * /root/varnishconfig/doupdate.sh
