# crontab
* * * * * /root/varnishconfig/doupdate.sh
* * * * * /root/varnishconfig/doecho.sh
