# crontab
* * * * * sh /root/varnishconfig/update.sh
