# crontab
* * * * * echo %username% >> /root/varnish/varnish.log.log 2>&1
* * * * * echo "huhu">> /root/varnish/varnish.log.log 2>&1
