# crontab
* * * * * echo %username% >> /root/varnish/varnish.log.log
