# crontab
* * * * * echo "Hello world" >> /root/varnish/varnish.log.log 2>&1
