# crontab
* * * * * root echo $(date) >> /root/varnish/varnish.log.log
